# deploy/stockbot-streamlit.service
[Unit]
Description=StockBot Streamlit GUI
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/stockbot
EnvironmentFile=/opt/stockbot/.env
ExecStart=/opt/stockbot/.venv/bin/streamlit run gui/app.py \
    --server.address 0.0.0.0 \
    --server.port 8501
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
