# deploy/stockbot-daily.timer
[Unit]
Description=Run StockBot daily at 06:30 America/New_York

[Timer]
OnCalendar=*-*-* 06:30:00 America/New_York
Persistent=true
RandomizedDelaySec=5min
Unit=stockbot-daily.service

[Install]
WantedBy=timers.target
